# traffic.service
[Unit]
Description=Traffic Light Controller
After=network.target

[Service]
Type=simple
User=pi
Group=gpio
WorkingDirectory=/home/pi/Trafficlight-2.0
ExecStart=/usr/bin/python3 /home/pi/Trafficlight-2.0/traffic.py
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_NET_BIND_SERVICE
StandardOutput=journal
StandardError=journal
Environment=TRAFFIC_PORT=80
Environment=TRAFFIC_KEYHOLDER_PIN=1234
Environment=TRAFFIC_SECRET_KEY=change-me-in-production

[Install]
WantedBy=multi-user.target
Alias=trafficlight.service
